# infra/rc.d/bsdos_tunnel
#!/bin/sh
#
# bsdos-tunnel: Wayland wire proxy for foot pipeline (wl_shm → stream socket)
#
# PROVIDE: bsdos_tunnel
# REQUIRE: bsdos_cage
# BEFORE: bsdos_core
# KEYWORD: shutdown
#
# Add to /etc/rc.conf:
#   bsdos_tunnel_enable="YES"

. /etc/rc.subr

name="bsdos_tunnel"
rcvar="bsdos_tunnel_enable"

load_rc_config $name

: ${bsdos_tunnel_enable:="NO"}
: ${bsdos_tunnel_rundir:="/tmp/wayland-run"}
: ${bsdos_tunnel_compositor_sock:="/tmp/wayland-run/wayland-0"}
: ${bsdos_tunnel_ghost_sock:="/tmp/wayland-run/wayland-ghost-0"}
: ${bsdos_tunnel_stream_sock:="/tmp/wayland-run/wayland-stream.sock"}

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
procname="/usr/local/bin/wayland-tunnel"
start_cmd="${name}_start"
stop_cmd="${name}_stop"

bsdos_tunnel_wait_socket() {
    local i=0
    while [ $i -lt 30 ]; do
        [ -S "${bsdos_tunnel_ghost_sock}" ] && [ -S "${bsdos_tunnel_stream_sock}" ] && return 0
        sleep 0.5
        i=$((i + 1))
    done
    echo "${name}: timeout waiting for sockets" >&2
    return 1
}

bsdos_tunnel_kill_stale() {
    # wayland-tunnel cmdline is just the binary name — socket paths live in env vars only.
    # pkill -f pattern on env-only args never matches; use sockstat to find by socket dir.
    sockstat 2>/dev/null \
        | awk '/wayland-tu.*wayland-run/{print $3}' \
        | sort -u \
        | xargs -r kill -9 2>/dev/null || true
    rm -f "${bsdos_tunnel_ghost_sock}" "${bsdos_tunnel_ghost_sock}.lock" \
          "${bsdos_tunnel_stream_sock}" 2>/dev/null || true
}

bsdos_tunnel_start() {
    # Require cage compositor
    if [ ! -S "${bsdos_tunnel_compositor_sock}" ]; then
        echo "${name}: compositor socket ${bsdos_tunnel_compositor_sock} not found, skipping"
        return 1
    fi

    bsdos_tunnel_kill_stale
    sleep 0.5

    echo "Starting ${name}..."

    daemon -p ${pidfile} -r -o "${bsdos_tunnel_rundir}/wayland-tunnel.log" \
        env XDG_RUNTIME_DIR="${bsdos_tunnel_rundir}" \
            BSDOS_COMPOSITOR_SOCK="${bsdos_tunnel_compositor_sock}" \
            BSDOS_GHOST_SOCK="${bsdos_tunnel_ghost_sock}" \
            BSDOS_STREAM_SOCK="${bsdos_tunnel_stream_sock}" \
        /usr/local/bin/wayland-tunnel

    if bsdos_tunnel_wait_socket; then
        echo "${name} started"
    else
        echo "${name}: WARNING - sockets not created"
    fi
}

bsdos_tunnel_stop() {
    echo "Stopping ${name}..."
    if [ -f "${pidfile}" ]; then
        kill -TERM $(cat ${pidfile} 2>/dev/null) 2>/dev/null || true
        rm -f "${pidfile}"
    fi
    bsdos_tunnel_kill_stale
}

run_rc_command "$1"
